# deploy/texo-agent.service
[Unit]
Description=texo memory agent (claim-chain memory over Qwen Cloud)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/opt/texo/env
Environment=TEXO_AGENT_ROOT=/opt/texo/workspace
Environment=TEXO_AGENT_ADDR=0.0.0.0:8787
Environment=TEXO_UI_DIR=/opt/texo/ui/dist
ExecStart=/opt/texo/bin/texo serve
Restart=on-failure
RestartSec=2
# The journal is the durable state; give it a stable home.
WorkingDirectory=/opt/texo/workspace

[Install]
WantedBy=multi-user.target
